# server-conf/supervisor-uwsgi.conf
[program:uwsgi]
command = /usr/bin/uwsgi --ini /home/pi/pi_cam/server-config/django_project_uwsgi.ini

